# deploy/borger.service
[Unit]
Description=Borger Borgmatic Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/borger
ExecStart=/usr/bin/node /opt/borger/server.js
Environment=HOME=/root
Environment=APP_ADDR=:8080
Environment=APP_COLLECT_INTERVAL_SECONDS=600
Environment=BORGER_DB=/var/lib/borger/borger.db
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
