# deploy/borger.service
[Unit]
Description=Borger Borgmatic Dashboard
After=network.target

[Service]
Type=simple
User=borger
Group=borger
WorkingDirectory=/opt/borger
ExecStart=/usr/bin/node /opt/borger/server.js
Environment=APP_ADDR=:8080
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
